# devops/www/etc/cron.d/elastiq
* *    * * * root python /opt/elastiq/elastiq/manage.py deploy
